# counsel-day-server/systemd/luks-unlock.service
# /etc/systemd/system/luks-unlock.service
#
# Counsel.day · boot-time LUKS unlock + mount for the Postgres data volume.
#
# Runs BEFORE docker.service so any Postgres container that depends on
# the mounted volume cannot start before this unit succeeds.
#
# The passphrase is fetched from Infisical via the infisical CLI using
# a machine-identity service token loaded from /etc/environment.d/infisical.conf.
# If Infisical is unreachable, this unit fails and Postgres does not start ·
# Caddy will continue serving a maintenance page from a static volume.
#
# Install with:
#   sudo cp luks-unlock.sh /usr/local/sbin/luks-unlock.sh
#   sudo chmod 700 /usr/local/sbin/luks-unlock.sh
#   sudo cp luks-unlock.service /etc/systemd/system/luks-unlock.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable luks-unlock.service

[Unit]
Description=Counsel.day · LUKS unlock + mount for Postgres data volume
DefaultDependencies=no
After=network-online.target
Wants=network-online.target
Before=docker.service
RequiresMountsFor=/var/lib/postgres-data

# We need network up so the Infisical CLI can reach vault.counsel.day.
# We use network-online.target rather than network.target because the
# former waits for actual connectivity, not just an interface.

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/environment.d/infisical.conf
ExecStart=/usr/local/sbin/luks-unlock.sh
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal

# Hardening: this unit needs cryptsetup and mount, but nothing else.
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/postgres-data /run
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
