# counsel-day-app/ops/counsel-day-backup-verify.service
[Unit]
Description=Counsel.day · Postgres backup verification (weekly test-restore)
After=postgresql.service counsel-day-backup.service
Wants=postgresql.service

[Service]
Type=oneshot
User=deploy
Group=deploy
WorkingDirectory=/opt/counsel-day-app
EnvironmentFile=/etc/counsel-day-app/env.local
Environment="BACKUP_DIR=/var/backups/counsel-day"
ExecStart=/opt/counsel-day-app/scripts/backup-verify.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=counsel-day-backup-verify
TimeoutStartSec=30min

# Hardening · the verify job creates and drops its own throwaway DB on
# the local Postgres cluster. It does not touch the production schema.
# It needs to read the backup directory and write a temporary file
# while gunzip pipes through psql.
#
# NoNewPrivileges is intentionally NOT set true · the script calls
# `sudo -u postgres psql` to elevate to the postgres role for the
# CREATE/DROP DATABASE statements. The remaining hardening lines
# below still narrow filesystem and capability access.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
PrivateTmp=true
